# etc/cron/nzgl.cron
0 1 * * * root yum --quiet -y upgrade &>/dev/null
0 2 * * * root /home/build/nzgl/bin/log_updates; /home/build/nzgl/bin/mail_report
