# cron/nzgl.cron
0 1 * * * root /home/build/nzgl/cron/yum_update; /home/build/nzgl/cron/log_updates; /home/build/nzgl/cron/mail_report
